# launchd/com.vmlxstation.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vmlxstation.daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/ABSOLUTE/PATH/TO/vmlx-station/scripts/dev_daemon.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/REPLACE_ME/Library/Application Support/vmlx-station/logs/daemon.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_ME/Library/Application Support/vmlx-station/logs/daemon.launchd.log</string>
</dict>
</plist>

